# packages/carnotzet+rootfs-pi4b/root/rootfs/etc/inittab
::sysinit:/etc/init.d/rcS
tty1::respawn:/sbin/getty -L tty1 115200 vt100 /bin/login