# packages/carnotzet+rootfs-pi4b/root/rootfs/etc/inittab
::sysinit:/etc/init.d/rcS
ttyO2::respawn:/sbin/getty -L ttyO2 115200 vt100 /bin/login